# .config/systemd/user/every-five-minutes.service
[Unit]
Description=Execute actions every 5 minutes

[Service]
ExecStart=/usr/local/bin/every-five-minutes
