# .config/systemd/user/wireplumber.service.d/override.conf
[Service]
Environment=XDG_STATE_HOME="/home/dusan/.local/share"
